# quadlets/fabric-smoke-combat.container
[Unit]
Description=Fabric verify 3/4 — combat wire parity vs Lean trace
After=fabric-smoke-loot.service
[Container]
Image=localhost/fabric-smoke
Volume=/home/ernest.lee/Documents/godot/bin:/godot/bin:ro,z
Volume=/home/ernest.lee/Documents/fabric-verify/assets:/assets:ro,z
Volume=/home/ernest.lee/Documents/fabric-verify/bin:/smoke:ro,z
Exec=/smoke/run-combat.sh
[Service]
Type=oneshot
RemainAfterExit=yes
[Install]
WantedBy=fabric-verify.target
