# quadlets/fabric-verify.target
[Unit]
Description=Fabric verification queue — Lean-grounded smokes under systemd
Wants=fabric-smoke-monado.service fabric-smoke-loot.service fabric-smoke-combat.service fabric-smoke-fourplayer.service fabric-smoke-loop.service fabric-smoke-mobile-shaders.service
