# quadlets/fabric-smoke-fourplayer.container
[Unit]
Description=Fabric verify 4/4 — four-player contention vs Lean golden
After=fabric-smoke-combat.service
[Container]
Image=localhost/fabric-smoke
Volume=/home/ernest.lee/Documents/godot/bin:/godot/bin:ro,z
Volume=/home/ernest.lee/Documents/fabric-verify/assets:/assets:ro,z
Volume=/home/ernest.lee/Documents/fabric-verify/bin:/smoke:ro,z
Exec=/smoke/run-fourplayer.sh
[Service]
Type=oneshot
RemainAfterExit=yes
[Install]
WantedBy=fabric-verify.target
